# database/graphite/templates/upstart-graphite.conf
# graphite
#

description "graphite server"

start on virtual-filesystems
start on runlevel [23]
stop on runlevel [23]

exec /opt/graphite/bin/carbon-cache.py start

respawn