# database/redis/templates/upstart-redis.conf
# redis
#

description "redis server"

start on virtual-filesystems
start on runlevel [23]
stop on runlevel [23]

exec sudo -u redis -g adm /usr/bin/redis-server /etc/redis/redis.conf

respawn